# files/conf.d/asterisk.conf
[program:asterisk]
command=/usr/sbin/asterisk -f -C /etc/asterisk/asterisk.conf -U asterisk -G asterisk
autostart=true
autorestart=true
startretries=10
